# templates/default/carbon-cache.conf.erb
# THIS FILE HAS BEEN GENERATED BY CHEF
# ANY MANUAL MODIFICATIONS WILL BE OVERWRITTEN

description "carbon-cache v<%= node['graphite']['version'] %>"

start on runlevel [2345]
stop on runlevel [!2345]

env CMD=/opt/graphite/bin/carbon-cache.py
env PIDFILE=/opt/graphite/storage/carbon-cache.pid

# uncomment this for debugging
# console output

pre-start script
  $CMD --pidfile=$PIDFILE start
end script

post-stop script
  $CMD --pidfile=$PIDFILE stop

  # don't consider this service stopped until the PID gets removed
  while [ -e $PIDFILE ]; do sleep 1 ; done
  exit 0
end script
