# templates/default/carbon-cache.conf.erb
description	"Carbon server"
author "Chef"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

umask 022
expect daemon
respawn

pre-start script
    test -d /opt/graphite || { stop; exit 0; }
end script

chdir /opt/graphite

exec start-stop-daemon --oknodo --chdir /opt/graphite --user <%= @user %> --chuid <%= @user %> --pidfile /opt/graphite/storage/carbon-cache-a.pid --name carbon-cache --startas /opt/graphite/bin/carbon-cache.py --start start
